# conf/supervisord/supervisor.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command=sh /startPhpFpm.sh
autostart=true
autorestart=true

[program:nginx]
command=sh /startNginx.sh
autostart=true
autorestart=true
